# media-drop.service
[Unit]
Description=Media Drop Downloader
After=network.target

[Service]
Type=simple
User=matt
WorkingDirectory=/srv/media-drop/app
ExecStart=/usr/bin/node server.js
Restart=on-failure
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=MEDIA_DROP_STORAGE_ROOT=/srv/media-drop
Environment=MEDIA_DROP_DB_PATH=/srv/media-drop/db/media-drop.sqlite
Environment=MEDIA_DROP_ADMIN_PASSWORD=your-secure-password
Environment=MEDIA_DROP_MAX_CONCURRENT_JOBS=3
Environment=MEDIA_DROP_MAX_FILE_SIZE_GB=2

[Install]
WantedBy=multi-user.target
